# airootfs/etc/profile.d/calamares.sh
#!/usr/bin/env bash
chmod +755 /airootfs/etc/skel/Desktop/calamares.desktop
